# ollama.service
[Unit]
Description=Ollama AI Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ollama serve
Restart=always
User=gabriel

[Install]
WantedBy=multi-user.target